# eink-display.service
[Unit]
Description=E-Ink Display Service
After=network.target

[Service]
Type=simple
User=asperkins
WorkingDirectory=/home/pi/waveshare-pinfo
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/pi/waveshare-pinfo/.venv/bin/python /home/pi/waveshare-pinfo/eink_display.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
