# eink-display.service
[Unit]
Description=E-Ink Display Service
After=network.target dev-spidev0.0.device

[Service]
Type=simple
User=asperkins
Group=asperkins
SupplementaryGroups=gpio,spi
WorkingDirectory=/opt/waveshare-pinfo
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 2
ExecStart=/opt/waveshare-pinfo/.venv/bin/python /opt/waveshare-pinfo/eink_display.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
