# photo-display.service
[Unit]
Description=Google Photos E-ink Display Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/waveshare-pinfo
Environment=PYTHONPATH=/opt/waveshare-pinfo
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/waveshare-pinfo/.venv/bin/python /opt/waveshare-pinfo/photo_display_service.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Ensure network connectivity before starting
ExecStartPre=/bin/bash -c 'until ping -c1 google.com; do sleep 1; done'

# Kill any existing processes on restart
ExecStartPre=/bin/bash -c 'pkill -f photo_display_service.py || true'

# Environment variables (optional - can also use .env file)
# Environment=GOOGLE_PHOTOS_ALBUM_ID=your_album_id
# Environment=PHOTO_ROTATION_INTERVAL_MINUTES=30

[Install]
WantedBy=multi-user.target
